# docker/cron/crontab
* * * * * cd /var/www/html && /usr/local/bin/php artisan schedule:run >> /var/www/html/docker/cron/cron_log.log 2>&1
